# systemd-timers/backup.service
[Unit]
Description=Full System Backup

[Service]
Type=simple
ExecStart=/home/thegner/tools/backup/backup.sh
