# systemd-timers/backup.timer
[Unit]
Description=Regular Backups

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
